# deploy/openfang.service
[Unit]
Description=OpenParlant Agent OS Daemon
Documentation=https://github.com/openparlant-ai/openparlant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openparlant
Group=openparlant
ExecStart=/usr/local/bin/openparlant start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
Environment=OPENFANG_HOME=/var/lib/openparlant
EnvironmentFile=-/etc/openparlant/env
WorkingDirectory=/var/lib/openparlant

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openparlant
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
