# deploy/openfang.service
[Unit]
Description=OpenParlant Agent OS Daemon
Documentation=https://github.com/silicrew-ai/silicrew
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=silicrew
Group=silicrew
ExecStart=/usr/local/bin/silicrew start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
Environment=OPENFANG_HOME=/var/lib/silicrew
EnvironmentFile=-/etc/silicrew/env
WorkingDirectory=/var/lib/silicrew

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/silicrew
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
